# supervisord.conf
[unix_http_server]
file=/tmp/supervisor.sock
username=mybnbaid
password=password

[inet_http_server]
port=localhost:9999
username=mybnbaid
password=password

[supervisord]
logfile=%(here)s/logs/supervisor.log
logfile_maxbytes=50MB
logfile_backups=5
loglevel=info
pidfile=%(here)s/logs/supervisor.pid
nodaemon=false
minfds=4096
minprocs=128
environment=PYTHONPATH=%(here)s

[rpcinterface:supervisor]
supervisor.rpcinterface_factory=supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[group:devices]
programs=users.queue, mail.streamer 

[group:gateway]
programs=gateway.service

[group:services]
programs=users.service, mail.service

[program:users.queue]
command=rock.device -s users -d queue -f tcp://*:5001 -b tcp://*:5000
directory=%(here)s
process_name=users.queue
numprocs=1
priority=1  ; start first
autostart=true
autorestart=true
startsecs=10
stopwaitsecs=10
killasgroup=true
stdout_logfile=%(here)s/logs/users.queue.log
stderr_logfile=%(here)s/logs/users.queue.log

[program:mail.streamer]
command=rock.device -s mail -d streamer -f tcp://*:6001 -b tcp://*:6000
directory=%(here)s
process_name=mail.streamer
numprocs=1
priority=1  ; start first
autostart=true
autorestart=true
startsecs=10
stopwaitsecs=10
killasgroup=true
stdout_logfile=%(here)s/logs/mail.streamer.log
stderr_logfile=%(here)s/logs/mail.streamer.log

[program:users.service]
command=mybnbaid.users
directory=%(here)s
process_name=users.service-[%(process_num)02d]
numprocs=2
priority=10  ; start third
autostart=true
autorestart=true
startsecs=10
stopwaitsecs=10
killasgroup=true
stdout_logfile=%(here)s/logs/users.service.log
stderr_logfile=%(here)s/logs/users.service.log

[program:mail.service]
command=mybnbaid.mail
directory=%(here)s
process_name=mail.service-[%(process_num)02d]
numprocs=2
priority=20  ; start third
autostart=true
autorestart=true
startsecs=10
stopwaitsecs=10
killasgroup=true
stdout_logfile=%(here)s/logs/mail.service.log
stderr_logfile=%(here)s/logs/mail.service.log

[program:gateway.service]
command=mybnbaid.gateway
directory=%(here)s
process_name=gateway.service-[%(process_num)02d]
numprocs=1
priority=30  ; start fourth
autostart=true
autorestart=true
startsecs=10
stopwaitsecs=10
killasgroup=true
stdout_logfile=%(here)s/logs/gateway.service.log
stderr_logfile=%(here)s/logs/gateway.service.log
